# picker/picker_camera_still_startup.service
[Unit]
Description=Picker Camera Still Startup Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=gregm
WorkingDirectory=/home/gregm/hardware_exercises
Environment=PYTHONPATH=/home/gregm/hardware_exercises
ExecStart=/home/gregm/hardware_exercises/.venv/bin/python3 picker/run_picker.py --generation-mode img2img --verbose --rotary
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target